# quraan-werd-bot.service
[Unit]
Description=Quraan Werd Tracker Bot
After=network.target

[Service]
Type=simple
User=your_username
WorkingDirectory=/path/to/bot/directory
ExecStart=/usr/bin/python3 /path/to/bot/directory/bot.py
Restart=on-failure
RestartSec=10
StandardOutput=syslog
StandardError=syslog
SyslogIdentifier=quraan-werd-bot
Environment=PYTHONUNBUFFERED=1

[Install]
WantedBy=multi-user.target 